# deploy/quadlets/podman-exporter.container
# prometheus-podman-exporter — per-container metrics (N6.7). Reads the rootless
# podman socket (enable it once: `systemctl --user enable --now podman.socket`).
# Not published; scraped on blog.network.
[Unit]
Description=podman exporter (container metrics)
Requires=blog-network.service podman.socket
After=blog-network.service podman.socket

[Container]
Image=quay.io/navidys/prometheus-podman-exporter:v1.12.0
ContainerName=podman-exporter
Network=blog.network
Environment=CONTAINER_HOST=unix:///run/podman/podman.sock
# %t is XDG_RUNTIME_DIR; the rootless podman socket lives there. Mount it
# READ-WRITE (connecting to a unix socket needs write), and run as container-root
# — in rootless podman that maps to the blog user that owns the socket, so the
# exporter can reach it (a non-root container uid maps to a subuid that cannot).
Volume=%t/podman/podman.sock:/run/podman/podman.sock
User=root
SecurityLabelDisable=true

[Service]
Restart=always

[Install]
WantedBy=default.target
